# scripts/jarvis-hardware-tier.service
# jarvis-hardware-tier.service
# Drop into /etc/systemd/system/ and enable:
#   sudo cp jarvis-hardware-tier.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jarvis-hardware-tier.service
#
# Writes /run/jarvis/tier.json so the llm-router (and any other consumer)
# can decide between local-GPU and cloud-only mode. Required because /run
# is tmpfs and the file disappears on reboot; this unit recreates it.

[Unit]
Description=Write /run/jarvis/tier.json so the llm-router knows GPU vs CPU
After=network-online.target
Before=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/jarvis-network/scripts/detect-hardware-tier.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
